# crontabs/root
MAILTO="dummy-mailto"
* * * * * /bin/not_found
